# com.benjamin.musicdashboardsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.benjamin.musicdashboardsync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/benjaminbellman/Music Dashboard/run_sync.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/benjaminbellman/Music Dashboard</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/benjaminbellman/Music Dashboard/logs/sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/benjaminbellman/Music Dashboard/logs/sync.log</string>
</dict>
</plist>
